# api-server/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:gunicorn]
command=gunicorn --config gunicorn_config.py RunServer:app

[program:nginx]
command=nginx -g "daemon off;"
